# root/etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
Execstart=
ExecStart=-/usr/bin/agetty --autologin wael --noclear %I $TERM